# packaging/sshtmd.service
[Unit]
Description=SSH Tunnel Manager Daemon

[Service]
Type=simple
ExecStart=/usr/bin/sshtmd

[Install]
WantedBy=default.target
